# nfs.automount.in
[Unit]
Description=Automount NFS Secure Offline Updates

[Automount]
Where=/nfs

[Install]
WantedBy=multi-user.target
